# swa.service
[Unit]
Description=NOAA Space Weather Alerts
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
Environment="DISPLAY=:0"
ExecStart=screen -DmS swa /bin/bash -c '/usr/local/bin/swa -ar; exec bash'
Restart=on-failure

[Install]
WantedBy=default.target
